# tarrot-server.service
[Unit]
Description=Tarrot Leerplatform Server
After=network.target

[Service]
Type=simple
User=fb-techlabs3
WorkingDirectory=/home/fb-techlabs3/rescue/tarrot
ExecStart=/home/fb-techlabs3/.hermes/node/bin/node node_modules/.bin/tsx server.ts
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
EnvironmentFile=/home/fb-techlabs3/rescue/tarrot/.env

[Install]
WantedBy=multi-user.target
